# overlay/etc/s6-overlay/s6-rc.d/apache/finish
#!/usr/bin/nu

use bf
use bf-apache
use bf-s6
bf env load -x apache

# Ensure Apache is closed gracefully and use base executable to terminate all services.
# If we are restarting, do nothing.
def main [...args] {
    # if Apache is marked as being restarted, do nothing
    if (bf-apache ctl is_restarting) {
        bf write debug "Apache is restarting, do not terminate container."
        return
    }

    # ensure Apache is stopped gracefully
    bf-apache ctl stop

    # finish and terminate container
    bf-s6 svc finish --terminate
}
