# overlay/etc/s6-overlay/s6-rc.d/apache/run
#!/bin/withenv sh

set -euo pipefail


#======================================================================================================================
# Run preflight checks
#======================================================================================================================

nu -c "use bf-apache run ; run preflight"


#======================================================================================================================
# Start Apache as a foreground service
#
# Flags:
#   -D  set daemon mode
#   -f  set configuration file path
#======================================================================================================================

apachectl start -D FOREGROUND -f "${BF_APACHE_CONF}"
